# glslave/etc/cron.d/glslave.cron
15 * * * * root /usr/bin/CronRefreshNetwork.py
