# glmaster/etc/cron.d/glmaster.cron
15 * * * * root /usr/bin/CronOpenPort.py
