# glmaster/etc/cron.d/glslave.cron
15 * * * * root /usr/bin/CronOpenPort.py
15 * * * * root /usr/bin/CronRefreshNetwork.py
