# domomaster/etc/cron.d/domomaster
15 * * * * root /etc/domoleaf/cron/CronOpenPort.py
30 * * * * root /etc/domoleaf/cron/CronLogsCalc.py
