# domoslave/etc/cron.d/domoslave
15 * * * * root /etc/domoleaf/cron/CronRefreshNetwork.py
*/15 * * * * root /etc/domoleaf/cron/CronSendAlive.py
